# pvheating.service
[Unit]
Description=Check and start PVHeating screen session
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/check_pvheating.sh
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
