# pvheating/DEBIAN/postinst
#!/bin/bash
set -e
systemctl daemon-reload
systemctl enable pvheating.service
systemctl start pvheating.service
