# pvheating/DEBIAN/prerm
#!/bin/bash
systemctl stop pvheating.service
systemctl disable pvheating.service
